# deploy/marketing-manager.service
[Unit]
Description=Marketing Manager Telegram Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/marketing-manager
EnvironmentFile=/home/ubuntu/marketing-manager/.env
ExecStart=/home/ubuntu/marketing-manager/venv/bin/python bot_telegram.py
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/marketing-manager/logs/bot.log
StandardError=append:/home/ubuntu/marketing-manager/logs/bot.log

[Install]
WantedBy=multi-user.target
